# webspider.service
[Unit]
Description=Web Spider server
After=online.target

[Service]
Type=simple
LimitCORE=infinite
ExecStart=/home/webspider/bin/webspider
WorkingDirectory=/home/webspider/www
Restart=always
TimeoutSec=5

[Install]
WantedBy=multi-user.target
